# dist/upstart/sa-monitoring-client.conf
# Upstart: /etc/init/sa-monitoring-client.conf
description "SA Monitoring Client"

start on filesystem
stop on shutdown

respawn

pre-start script
    test -x /usr/bin/sa-monitoring-client || { stop; exit 0; }
end script

exec /usr/bin/sa-monitoring-client
