# dist/upstart/sa-monitoring-server.conf
# Upstart: /etc/init/sa-monitoring-server.conf
description "SA Monitoring Server"

start on filesystem
stop on shutdown

respawn

pre-start script
    test -x /usr/bin/sa-monitoring-server || { stop; exit 0; }
end script

exec /usr/bin/sa-monitoring-server
