# dist/chef/cookbooks/sensu/templates/default/init/sensu-service.conf.erb
# Upstart: /etc/init/sensu-<%= @service %>.conf
# Managed by OpsChef

description "sensu <%= @service %>"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on [!12345]

respawn

pre-start script
  test -x <%= node.sensu.bin_path %>/sensu-<%= @service %> || { stop; exit 0; }
end script

exec setuidgid <%= node.sensu.user %> <%= node.sensu.bin_path %>/sensu-<%= @service %><% if @options %> <%= @options %><% end %>
