# dist/chef/cookbooks/sensu/templates/default/init/sensu-service.conf.erb
# Upstart: /etc/init/sensu-<%= @service %>.conf
# Managed by OpsChef

description "sensu <%= @service %>"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on [!12345]

respawn

chdir <%= node.sensu.bin_path %>
env PATH=<%= node.sensu.directory %>/plugins:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

pre-start script
  test -x ./sensu-<%= @service %> || { stop; exit 0; }
end script

exec setuidgid <%= node.sensu.user %> ./sensu-<%= @service %><% if @options %> <%= @options %><% end %>
